# root.cron
0 * * * *	/usr/local/rightscale-plugin-costs/rl-overage-plugin-cost.py
